# xinit.service
[Unit]
Description=~/.xinitrc

[Service]
ExecStart=/bin/bash -l .xinitrc

[Install]
WantedBy=xinit.target
