# xinit.target
[Unit]
Description=xinit session
Wants=xinit.service
AllowIsolate=yes

[Install]
Alias=default.target
